# factorio-init@.service
[Unit]
Description=initialize factorio instance %i

[Service]
Type=oneshot
ExecStart=/usr/local/factorio/init.sh %i
