# factorio-update.timer
[Unit]
Description=Run Factorio server binary auto-update hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
